# conf/sys_monitor.service
[Unit]
Description = System Monitor service

[Service]
ExecStart = /usr/local/bin/sys_monitor

[Install]
WantedBy = multi-user.target 
